# roles/rpi-argonone/files/argononed-init
#!/sbin/openrc-run

name="argononed"
command=/usr/local/bin/argononed.py
command_background=yes
pidfile=/run/$RC_SVCNAME.pid
required_files=/etc/argononed.conf
