# roles/alpine-physical/files/argononed-init
#!/sbin/openrc-run

name="Argononed"
command=/usr/local/bin/argononed.py
command_background=yes
pidfile=/run/$RC_SVCNAME.pid
required_files=/etc/argononed.conf
