# roles/alpine-physical/templates/inittab
# /etc/inittab from @juliekoubova/machines/roles/alpine-physical

::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default

# Set up a couple of mingetty's
{% if alpine_release == "v3.20" %}
{%  set mingetty = "/sbin/mingetty" %}
{% else %}
{%  set mingetty = "/usr/sbin/mingetty" %}
{% endif %}
{% if autologin %}
tty1::respawn:/bin/login -f {{ autologin }}
{% else %}
tty1::respawn:{{ mingetty }} tty1
{% endif %}
tty2::respawn:{{ mingetty }} tty2
tty3::respawn:{{ mingetty }} tty3
tty4::respawn:{{ mingetty }} tty4
tty5::respawn:{{ mingetty }} tty5
tty6::respawn:{{ mingetty }} tty6

# Put a getty on the serial port
#ttyS0::respawn:/sbin/getty -L ttyS0 115200 vt100

# Stuff to do for the 3-finger salute
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
::shutdown:/sbin/openrc shutdown


